# home/.crontab
MAILTO=""
* 23 * * * /usr/local/bin/brew update &> /dev/null
